# kafka/supervisor/kafka-connect.ini
[program:kafka-connect]
command=sh /usr/bin/start-kafka-connect.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0